# backend/crontab
*/30 * * * * /usr/local/bin/cleanup.sh


